# hostcontrold.service
[Unit]
Description=Host Control Daemon

[Service]
ExecStart=/usr/sbin/hostcontrold
User=root

[Install]
WantedBy=multi-user.target